# supervisord.conf
[supervisord]
nodaemon=true
user=root

[program:xvfb]
# Comando crítico: Se declara explícitamente una resolución de 1920x1080
# para evitar que el renderizado caiga en la resolución 1x1 que los WAF bloquean.[3, 17]
command=/usr/bin/Xvfb :1 -screen 0 1920x1080x24 -nolisten tcp
autorestart=true
priority=100

[program:fluxbox]
# Lanza el gestor de ventanas atado a la pantalla virtual
command=/usr/bin/fluxbox -display :1
environment=DISPLAY=":1"
autorestart=true
priority=200

[program:x11vnc]
# Se eliminó la bandera -bg para que Supervisor pueda monitorearlo en primer plano
command=/usr/bin/x11vnc -display :1 -rfbauth /root/.vnc/passwd -forever -shared
autorestart=true
priority=300

[program:novnc]
# Websockify escucha el puerto interno de VNC (5900) y lo expone en HTML5 por el puerto 8080 
command=/usr/bin/websockify --web=/usr/share/novnc 8080 localhost:5900
autorestart=true
priority=400

[program:camoufox_hex_app]
# Ejecuta la capa de aplicación principal
command=/opt/venv/bin/python /app/main.py
environment=DISPLAY=":1",PYTHONUNBUFFERED="1"
autorestart=false
priority=500
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0